# .huggingface/supervisord.conf
[supervisord]
nodaemon=true
logfile=/app/logs/supervisord.log
pidfile=/tmp/supervisord.pid
user=root

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/app/logs/nginx-stdout.log
stderr_logfile=/app/logs/nginx-stderr.log
priority=10

[program:fastapi]
command=uvicorn api.main:app --host 0.0.0.0 --port 8000 --log-level info --proxy-headers
directory=/app
autostart=true
autorestart=true
stdout_logfile=/app/logs/fastapi-stdout.log
stderr_logfile=/app/logs/fastapi-stderr.log
environment=PYTHONUNBUFFERED="1",HF_SPACES="1"
priority=20

[program:docusaurus]
command=serve /app/static/docs -l 3000 -n
directory=/app
autostart=true
autorestart=true
stdout_logfile=/app/logs/docusaurus-stdout.log
stderr_logfile=/app/logs/docusaurus-stderr.log
priority=30

[program:frontend]
command=serve /app/static/frontend -l 5173 -n
directory=/app
autostart=true
autorestart=true
stdout_logfile=/app/logs/frontend-stdout.log
stderr_logfile=/app/logs/frontend-stderr.log
priority=30
